# scripts/schedules/com.zeego.resourcematch-seo-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.zeego.resourcematch-seo-bot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/zeego/Engineering/resourcematch/scripts/venv/bin/python</string>
        <string>-m</string>
        <string>scripts.seo.bot</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/zeego/Engineering/resourcematch</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>/Users/zeego/Engineering/resourcematch/scripts/logs/seo_bot_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/zeego/Engineering/resourcematch/scripts/logs/seo_bot_stderr.log</string>
</dict>
</plist>
